# .config/systemd/user/update-background.service
[Unit]
Description=Update Background
 
[Service]
Type=simple
ExecStart=/home/weixing/.config/systemd/user/update-background.sh